# reports/inflight/job-0314/titiler.service
[Unit]
Description=TiTiler COG tile server (GRACE-2 AWS raster rendering, job-0290; hardened job-0314)
After=network.target

[Service]
Environment=AWS_REGION=us-west-2
Environment=AWS_DEFAULT_REGION=us-west-2
Environment=GDAL_DISABLE_READDIR_ON_OPEN=EMPTY_DIR
Environment=CPL_VSIL_CURL_ALLOWED_EXTENSIONS=.tif,.tiff,.TIF
Environment=GDAL_HTTP_MERGE_CONSECUTIVE_RANGES=YES
# job-0314 RESILIENCE — the 2026-06-16 wedge root cause: a stalled /vsis3 read
# blocked a uvicorn worker forever (no timeout), threadpool exhausted across
# both workers, all tile requests then hung. These bound every S3 read so a
# stalled read fails fast and the worker thread frees instead of wedging.
Environment=GDAL_HTTP_TIMEOUT=30
Environment=GDAL_HTTP_CONNECTTIMEOUT=10
Environment=GDAL_HTTP_MAX_RETRY=2
Environment=GDAL_HTTP_RETRY_DELAY=1
Environment=VSI_CACHE=TRUE
Environment=VSI_CACHE_SIZE=26214400
Environment=GDAL_CACHEMAX=512
Environment=CPL_VSIL_CURL_CACHE_SIZE=209715200
# CORS scoped to the deployed site origin + local dev (NOT wildcard). The
# CloudFront-served site fetches /cog tiles same-origin, so CORS is only for
# direct S3-website access. Unchanged by job-0314.
Environment=TITILER_API_CORS_ORIGINS=http://grace2-hazard-web-226996537797.s3-website-us-west-2.amazonaws.com,http://localhost:5173
# job-0314 — 2 -> 4 workers for headroom under multi-Case tile bursts.
ExecStart=/opt/titiler/venv/bin/uvicorn titiler.application.main:app --host 0.0.0.0 --port 8080 --workers 4
Restart=always
RestartSec=2
User=root

[Install]
WantedBy=multi-user.target
